# server/deploy/jiankong-api.service
# 企业版抖音监控 — API（uvicorn）
#
# 安装（在服务器上执行，路径按实际修改）：
#   sudo cp /home/ubuntu/jiankong/server/deploy/jiankong-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jiankong-api
#
# 常用：
#   sudo systemctl status jiankong-api
#   sudo systemctl restart jiankong-api   # 代码 git pull 后用这条，不要用 kill
#   sudo systemctl stop jiankong-api
#
# 更新代码后：cd /home/ubuntu/jiankong && git pull && sudo systemctl restart jiankong-api

[Unit]
Description=Jiankong FastAPI (uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/jiankong/server
Environment=PYTHONUNBUFFERED=1
# 抗风控（S+A+B）完整示例见同仓库：server/deploy/jiankong-api.service.d/override.conf.example
# 与可选长 Cookie 文件：server/deploy/douyin.env.example
# Environment=DOUYIN_USE_PLAYWRIGHT=1
# Environment=DOUYIN_COOKIE=...  或 EnvironmentFile=-/etc/jiankong/douyin.env
# Environment=DOUYIN_PROXY_POOL=...   Environment=DOUYIN_PLAYWRIGHT_PROXY=...
ExecStart=/home/ubuntu/jiankong/.venv-enterprise/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
# 进程异常退出时重启；正常 stop 不会自动拉起。若需崩溃/杀进程也必自动重建可改为 Restart=always
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
